# install/volmon.service
[Unit]
Description=VolMon Volume Monitoring and Control Daemon
After=pulseaudio.service pipewire.service

[Service]
Type=simple
ExecStart=%h/.local/bin/volmon-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
